# misc/systemd/system/deepin-shutdown-sound.service
[Unit]
Description=Deepin shutdown sound
Requisite=sound.target local-fs.target
After=sound.target local-fs.target
Before=systemd-logind.service

[Service]
Type=forking
ExecStart=/usr/bin/true
ExecStop=/usr/lib/deepin-api/deepin-shutdown-sound
RemainAfterExit=yes
TimeoutStopSec=10s

[Install]
WantedBy=graphical.target
